# opg-mongodb-1604/docker/service/mongod/run
#!/bin/sh

gosu mongodb mongod -f /etc/mongod.conf
