# opg-ssmtp-1604/docker/etc/service/postfix-stats-qshape/run
#!/bin/sh

exec /usr/local/bin/postfix-stats-qshape
